# .config/systemd/user/picom.service
[Unit]
Description=Compositor for Xorg

[Service]
ExecStart=/usr/bin/picom -b

[Install]
WantedBy=default.target
